# launch-agents/com.example.calendar-work-color-daily.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- TEMPLATE. See the maintain-daily template for placeholder meanings.
     Work calendar: color-only (no emoji), meetings classified via the career default. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>__LABEL__</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>set -o pipefail; export HOME=__HOME__; PY="__PY__"; CFG="__INSTALL__/gcal/config.py"; CALS=$("$PY" "$CFG" --calendar-ids --account work); TOK=$("$PY" "$CFG" --token --account work); "$PY" "__INSTALL__/gcal/calendar_maintenance.py" --calendars "$CALS" --token "$TOK" --window-days 120 --restyle-all --emoji-off --default-category career --apply</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>6</integer>
		<key>Minute</key>
		<integer>10</integer>
	</dict>
</dict>
</plist>
