# launch-agents/com.example.calendar-learn-weekly.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- TEMPLATE. Replace the placeholders, then install (see README "Automation"):
       __LABEL__     e.g. com.you.calendar-learn-weekly
       __HOME__      your home dir, e.g. /Users/you
       __INSTALL__   where you cloned this repo, e.g. /Users/you/calendar-colors
     Weekly dynamic-taxonomy loop. Runs calendar-learn-weekly.sh (fill in its own
     placeholders first). REQUIRES an LLM agent CLI (the script uses `claude -p`). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>__LABEL__</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>__INSTALL__/launch-agents/calendar-learn-weekly.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>0</integer>
		<key>Hour</key>
		<integer>6</integer>
		<key>Minute</key>
		<integer>30</integer>
	</dict>
</dict>
</plist>
